# deployment/genred_english.service
[Unit]
Description=Genred English FastAPI Application
After=docker.service
Requires=docker.service

[Service]
Type=simple
Restart=always
RestartSec=10
TimeoutStartSec=0

# 環境変数ファイルを読み込む
EnvironmentFile=/opt/genred_english/.env

# 既存のコンテナを停止・削除
ExecStartPre=-/usr/bin/docker stop genred_english
ExecStartPre=-/usr/bin/docker rm genred_english

# コンテナを起動
ExecStart=/usr/bin/docker run --rm \
  --name genred_english \
  --network host \
  -e DATABASE_URL=${DATABASE_URL} \
  -e TRY_SOUND_BUCKET=${TRY_SOUND_BUCKET} \
  -e REGION_NAME=${REGION_NAME} \
  -e AWS_ACCESS_KEY_ID=${AWS_ACCESS_KEY_ID} \
  -e AWS_SECRET_ACCESS_KEY=${AWS_SECRET_ACCESS_KEY} \
  -e PYTHONUNBUFFERED=1 \
  genred_english:latest

# コンテナを停止
ExecStop=-/usr/bin/docker stop genred_english

[Install]
WantedBy=multi-user.target
